# src/supervisord.conf
# ──────────────────────────────────────────────────────────────────────────────
# supervisord.conf
# ──────────────────────────────────────────────────────────────────────────────
# IMPORTANT: supervisord requires every value to be on a SINGLE LINE.
# Multi-line command= blocks cause parse errors. All commands here are
# intentionally kept on one line each.
#
# Process startup order (controlled by priority):
#   100  postgres  — start first
#   200  fastapi   — start after postgres is ready
#   300  gradio    — start after fastapi is ready
# ──────────────────────────────────────────────────────────────────────────────

[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
loglevel=info

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

# ── PostgreSQL ─────────────────────────────────────────────────────────────────
[program:postgres]
command=/bin/bash -c "exec gosu postgres /usr/lib/postgresql/$(ls /usr/lib/postgresql/)/bin/postgres -D /data/pgdata -c listen_addresses=localhost -c log_destination=stderr"
priority=100
autostart=true
autorestart=true
startretries=3
startsecs=5
user=root
stdout_logfile=/tmp/postgres.log
stderr_logfile=/tmp/postgres.log

# ── FastAPI ────────────────────────────────────────────────────────────────────
[program:fastapi]
command=uvicorn main:app --host 0.0.0.0 --port 8000 --workers 1 --log-level info
directory=/app/backend
priority=200
autostart=true
autorestart=true
startretries=5
startsecs=10
environment=DATABASE_URL="%(ENV_DATABASE_URL)s"
stdout_logfile=/tmp/fastapi.log
stderr_logfile=/tmp/fastapi.log

# ── Gradio ─────────────────────────────────────────────────────────────────────
[program:gradio]
command=python app.py
directory=/app/frontend
priority=300
autostart=true
autorestart=true
startretries=5
startsecs=20
environment=BACKEND_URL="%(ENV_BACKEND_URL)s"
stdout_logfile=/tmp/gradio.log
stderr_logfile=/tmp/gradio.log